# deploy/logdbd.service
[Unit]
Description=logdbd — append-only audit log database
Documentation=https://github.com/lv-agent/logdb
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/logdbd --config /etc/logdbd/logdbd.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitMEMLOCK=infinity

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/logdbd
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
